# deploy/today_tactics.service
[Unit]
Description=Today Tactics — K리그 전술 분석
After=network.target

[Service]
User=rocky
Group=rocky
WorkingDirectory=/opt/today_tactics
Environment="PATH=/opt/today_tactics/venv/bin"
Environment="FLASK_DEBUG=0"
# FLASK_SECRET_KEY 미설정 시 재기동마다 새 키 생성 → 모든 로그인 세션 무효화.
# 운영 서버에서 반드시 고정 키를 주입할 것:
#   sudo systemctl edit today_tactics  →  Environment="FLASK_SECRET_KEY=<32바이트 이상 랜덤>" 추가
# Environment="FLASK_SECRET_KEY=<설정 필요>"

# CPU 우선순위 — 같은 서버의 다른 프로세스(bunting 봇) 폭주 시
# today_tactics가 사용자 응답 우선으로 CPU 받도록 비례 가중치 + nice 조정
# cgroup v2의 CPUWeight 기본값 100, 300은 약 3배 가중. Nice -5는 priority 약간 상향.
CPUWeight=300
Nice=-5
ExecStart=/opt/today_tactics/venv/bin/gunicorn \
    --workers 3 \
    --worker-class gthread \
    --threads 2 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --keep-alive 5 \
    --access-logfile /var/log/today_tactics/access.log \
    --error-logfile /var/log/today_tactics/error.log \
    main:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
